# templates/ood.service
[Unit]
Description=ood process manager
After=network.target

[Service]
ExecStart=%NODE% %OOD%

[Install]
WantedBy=multi-user.target

